# com.stockanalysis.batchupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- ジョブのラベル（一意の識別子） -->
    <key>Label</key>
    <string>com.stockanalysis.batchupdate</string>
    
    <!-- 実行するプログラム -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/hiranotakahiro/Projects/個別株分析/run_batch_update.sh</string>
    </array>
    
    <!-- 作業ディレクトリ -->
    <key>WorkingDirectory</key>
    <string>/Users/hiranotakahiro/Projects/個別株分析</string>
    
    <!-- 実行スケジュール：毎日午前2時に実行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <!-- 標準出力のログファイル -->
    <key>StandardOutPath</key>
    <string>/Users/hiranotakahiro/Projects/個別株分析/logs/launchd.out.log</string>
    
    <!-- 標準エラー出力のログファイル -->
    <key>StandardErrorPath</key>
    <string>/Users/hiranotakahiro/Projects/個別株分析/logs/launchd.err.log</string>
    
    <!-- システム起動時に自動で実行しない（時間指定のみ） -->
    <key>RunAtLoad</key>
    <false/>
    
    <!-- プロセスが終了しても自動で再起動しない -->
    <key>KeepAlive</key>
    <false/>
    
    <!-- バックグラウンドで実行 -->
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
